# etc/crontab
# добавить к /etc/crontab

01 12 * * * nginx /root/allsky.py/mp4.day.py
