# etc/crontab
# добавить к /etc/crontab

01 12 * * * root /opt/allsky.py/day.sh
* * * * * root /opt/allsky.py/sensors.sh
